# debian/cron.daily
#! /bin/sh

# update databases for both files and packages.
/usr/sbin/update-dlocatedb -b
