# debian/cron.daily
#!/bin/sh

# update databases for both files and packages.
[ -x /usr/sbin/update-dlocatedb ] && /usr/sbin/update-dlocatedb -b
