# debian/dlocate.cron.daily
#!/bin/sh

PATH="$PATH:/usr/bin:/bin:/usr/sbin:/sbin"

# update databases for both files and packages.
if command -v update-dlocatedb >/dev/null; then
  /usr/sbin/update-dlocatedb -b
fi
