# debian/dlocate.cron.daily
#!/bin/sh

# update databases for both files and packages.
if [ -x /usr/sbin/update-dlocatedb ]; then
  /usr/sbin/update-dlocatedb -b
fi
